# launchd/com.youssef.agent-island.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.youssef.agent-island</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/youssef/personal/agent-island/.build/release/AgentIsland</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/youssef/Library/Logs/AgentIsland.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/youssef/Library/Logs/AgentIsland.log</string>
</dict>
</plist>
